# ops/launchd/com.inkbaduk.ops-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 1시간마다 각 launchd 잡의 마지막 실행 timestamp를 검사하고 stale 시 알림을 보내는 watchdog. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.inkbaduk.ops-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/daegong/projects/baduk/ops/run-watchdog.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/daegong/projects/baduk/docs/ops/state/log/watchdog.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/daegong/projects/baduk/docs/ops/state/log/watchdog.err.log</string>
</dict>
</plist>
