# backend/deploy/com.baduk.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.baduk.api</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PROJECT_PATH__/backend/deploy/run_local_prod.sh</string>
  </array>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Interactive</string>

  <key>StandardOutPath</key>
  <string>__USER_HOME__/Library/Logs/baduk-api.log</string>
  <key>StandardErrorPath</key>
  <string>__USER_HOME__/Library/Logs/baduk-api.err</string>

  <key>WorkingDirectory</key>
  <string>__PROJECT_PATH__/backend</string>
</dict>
</plist>
